# docker/nodejs/node_supervisor.conf
[program:nodejs]
directory=/code
command=gulp prod
autostart=true
autorestart=true
redirect_stderr=true
redirect_stdout=true
stopsignal=INT
