# docker/nginx/nginx_supervisor.conf
[supervisord]
nodaemon       = true
logfile        = /var/log/supervisord/supervisord.log
pidfile        = /var/run/supervisord.pid 
childlogdir    = /var/log/supervisord/ 

[program:nginx]
command         = /usr/sbin/nginx
priority        = 999
autostart       = true
autorestart     = true
redirect_stderr = true
redirect_stdout = true
